# ansible/setup_monitoring/templates/node_exporter.service.j2
# node_exporter.service.j2
[Unit]
Description=Node Exporter

[Service]
User={{ node_exporter_user }}
ExecStart=/usr/local/bin/node_exporter
